# python_oracle/crontab
*/1 * * * * /bin/echo `date` >> /root/cron.log
